# contrib/crio-wipe/crio-wipe.service
[Unit]
Description=CRI-O Auto Update Script
Before=crio.service
RequiresMountsFor=/var/lib/containers

[Service]
ExecStart=/bin/bash /usr/local/libexec/crio/crio-wipe/crio-wipe.bash
Type=oneshot

[Install]
WantedBy=multi-user.target
